# whatever
#!/sbin/runscript
# Copyright 2013 Fanel
# Distributed under the terms of the GNU General Public License v2

depend() {
	after rogentoslive
}

start() {
	ebegin "Starting the..script"
	/usr/libexec/whatever.sh
	eend $?
}
